# root/home/haso/camviewer/service/buttons/run
#!/usr/bin/env bash

export PS_HOME=/home/haso/camviewer

cd "$PS_HOME"
python buttons.py

exit 0
